# systemd/slgpu.service
[Unit]
Description=slgpu — vLLM и/или SGLang через scripts/up.sh
Documentation=file:///opt/slgpu/README.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/slgpu
EnvironmentFile=-/opt/slgpu/.env

# Переопределите в drop-in (systemctl edit slgpu.service):
#   [Service]
#   Environment=SLGPU_MODE=both
#   Environment=MODEL=qwen3-30b-a3b
Environment=SLGPU_MODE=vllm
Environment=MODEL=

ExecStart=/usr/bin/bash -lc 'cd /opt/slgpu && exec ./scripts/up.sh "${SLGPU_MODE:-vllm}" ${MODEL:+-m "${MODEL}"}'
# Останавливаем только LLM-сервисы; Prometheus/Grafana остаются (см. README).
ExecStop=/usr/bin/bash -lc 'cd /opt/slgpu && docker compose stop vllm sglang || true'
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
